# ico-vic-carsim-wait-amb-ready.path
# Description=Path activation for ico-vic-carsim

[Path]
PathExists=/tmp/amd_ready
Unit=ico-vic-carsim.service

[Install]
WantedBy=weston.target
